# diabox/root/etc/profile.d/go.sh
export PATH=$PATH:~/go/bin:/go/bin:/usr/local/go/bin
